# .config/autostart/my-xmodmap.desktop
[Desktop Entry]
Name=MyXmodmap
Exec=/usr/bin/xmodmap /home/soimort/.Xmodmap
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
